# sunbird_rc/iac/roles/wireguard/templates/wireguard.service.j2
[Unit]
Description=wireguard docker compose unit
After=docker.service

[Service]
ExecStart=/usr/bin/docker compose -f {{ wireguard_compose_location }}/docker-compose.yml up
ExecStop=/usr/bin/docker compose -f {{ wireguard_compose_location }}/docker-compose.yml down

[Install]
WantedBy=docker.service